# templates/zfs-autobackup@.timer.j2
[Unit]
Description=zfs-autobackup

[Timer]
OnCalendar={{ zfs_autobackup_timer_schedule }}

[Install]
WantedBy=timers.target